# systemd/sgopherd.socket
[Socket]
ListenStream=70
Accept=yes

[Install]
WantedBy=sockets.target
